# service/udev/99-diskmanager.rules
#ENV{DEVTYPE}=="usb_device", ACTION=="add",RUN+="/usr/libexec/deepin-diskmanager/USBadd.sh"
#ENV{DEVTYPE}=="usb_device", ACTION=="remove", RUN+="/usr/libexec/deepin-diskmanager/USBremove.sh"

ENV{DEVTYPE}=="disk", ENV{ID_BUS}!="", ACTION=="add",RUN+="/usr/libexec/deepin-diskmanager/USBadd.sh"
ENV{DEVTYPE}=="disk", ENV{ID_BUS}!="", ACTION=="remove", RUN+="/usr/libexec/deepin-diskmanager/USBremove.sh"

